# scripts/ansible/roles/rpi/files/web_server.service
[Unit]
Description=Web Server Service
After=network.target

[Service]
ExecStart=/usr/local/bin/web_server
Restart=always
User=www-data
Group=www-data
Environment=PATH=/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/usr/local/bin

[Install]
WantedBy=multi-user.target

